# templates/etc/systemd/system/haproxy.service.d/override.conf.j2
[Service]
Environment="CONFIG=/etc/haproxy/" "PIDFILE=/run/haproxy.pid" "EXTRAOPTS=-S /run/haproxy-master.sock"
